# sys/etc/rc.local
#!/bin/sh

# /etc/rc.local -- rc.local for Artix Linux
#
# Enter your custom commands here.
